# otaready-app/org.webosarchive.otaready/device/otaready-daemon.conf
# /etc/event.d/otaready-daemon — upstart job for the Get-Ready-for-OTA status daemon.
description "Get Ready for OTA - status daemon"

start on started LunaSysMgr
stop on stopping LunaSysMgr

respawn
respawn limit 5 30

exec /usr/bin/otaready-daemon
